# castserial.service
[Unit]
Description=CastSerial
After=network.target

[Service]
ExecStart=/usr/local/sbin/castserial
WorkingDirectory=/home/pi-star
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi-star

[Install]
WantedBy=multi-user.target
